# platform/broadcom/sonic-platform-modules-accton/debian/sonic-platform-accton-as5822-54x.postinst
#Enable accton-as5822-platform services
depmod -a
systemctl enable as5822-platform-modules
systemctl start as5822-platform-modules
